# MB/Raspberry Antigas Automacoes Recebimento de ROM/Raspberry Automacoes/no var-www-html/AutomacaoGerdau/lib-systemd-system/gscs.service
[Unit]
Description=Executar_o_gscs_sva

[Service]
Type=simple
ExecStart=/home/pi/gscs.sh   

[Install]
WantedBy=multi-user.target

