# MB/Raspberry Antigas Automacoes Recebimento de ROM/Raspberry Automacoes/no var-www-html/AutomacaoGerdau/lib-systemd-system/saidas.service
[Unit]
Description=Executar_o_script_das_cancelas_saidas

[Service]
Type=simple
ExecStart=/home/pi/saidas.sh   

[Install]
WantedBy=multi-user.target

