# Automacao/Raspberry/Codigos e Dicas/Configurando Arquivos Iniciar automaticamente/crontab
@reboot sh /home/pi/automacao/launcher.sh >/home/pi/logs/cronlog 2>&1